# srv_test.service
[Unit]
Description=Service test

[Service]
ExecStart=/tools/repos/srv_test/srv_test.sh
User=nils
